# debian/clautod.postinst
#!/bin/sh
set -e

#DEBHELPER#

# Start the clautod service now that it's installed
systemctl start clautod

# Also setup the clautod service to run on system startup
systemctl enable clautod