# debian/clautod.service
[Unit]
Description=Consumables-Listing AUTOmation Daemon

[Service]
ExecStart=/home/jeremy/Projects/Clauto/clautod.sh /opt/venvs/clautod/bin/python3 /opt/venvs/clautod/lib/python3.5/site-packages/clautod/clautod.py

[Install]
WantedBy=multi-user.target