# debian/clautod.postrm
#!/bin/sh
set -e

#DEBHELPER#

# After uninstalling, don't run the clautod service anymore on system startup
systemctl disable clautod